# LaunchAgents/com.timr.login.plist
<?xml version="1.0" encoding="UTF-8"?>
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.timr.login</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/$(whoami)/timr/timr-login.sh</string>
    </array>
    <key>RunAtLoad</key><true/>
  </dict>
</plist>
